# launchd/com.songurov.discordbot.trans.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.songurov.discordbot.trans.daemon</string>

  <key>UserName</key>
  <string>songurov</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/songurov/DiscordBot &amp;&amp; /Users/songurov/.local/bin/trans start</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/songurov/DiscordBot</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StartInterval</key>
  <integer>120</integer>

  <key>StandardOutPath</key>
  <string>/Users/songurov/DiscordBot/.launchd-trans-daemon.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/songurov/DiscordBot/.launchd-trans-daemon.err.log</string>
</dict>
</plist>
